# server/release/datastation.service
[Unit]
Description=DataStation
After=network.target
StartLimitIntervalSec=0

[Service]
Type=simple
Restart=always
RestartSec=1
User=datastation
ExecStartPre=
WorkingDirectory=/usr/share/datastation
ExecStart=node /usr/share/datastation/build/server.js
ExecStartPost=
ExecStop=
ExecReload=

[Install]
WantedBy=multi-user.target